# rockiscope.service
[Unit]
Description=Rockiscope - Rockies Horoscope Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/rockiscope
EnvironmentFile=/home/pi/rockiscope/.env
ExecStart=/home/pi/rockiscope/rockiscope
Restart=always
RestartSec=30
StandardOutput=append:/home/pi/rockiscope/rockiscope.log
StandardError=append:/home/pi/rockiscope/rockiscope.log

[Install]
WantedBy=multi-user.target
